# scripts/systemd/hermes-auto-growth.service
[Unit]
Description=RA Hermes Non-Email Auto Growth Loop
After=network.target docker.service

[Service]
Type=oneshot
User=abyz-lab
WorkingDirectory=/home/abyz-lab/work/workspace-github/holee9/ra-hermes-multi-agent
EnvironmentFile=/home/abyz-lab/work/workspace-github/holee9/ra-hermes-multi-agent/scripts/.env
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=HONCHO_URL=http://localhost:8000
Environment=HONCHO_WS=work
Environment=AUTO_GROWTH_PENDING_SCOPE=ra
Environment=AUTO_GROWTH_MAX_PENDING=0
Environment=AUTO_GROWTH_CASES_PER_AGENT=1
Environment=AUTO_GROWTH_SOURCE_POOL=10
Environment=AUTO_GROWTH_MAX_CHUNKS_PER_CASE=1
Environment=AUTO_GROWTH_EXECUTE_CURRICULUM=weekly
Environment=PYTHONUNBUFFERED=1

ExecStart=/usr/bin/env bash scripts/auto-growth-runner.sh

StandardOutput=journal
StandardError=journal
SyslogIdentifier=hermes-auto-growth
TimeoutStartSec=21600

[Install]
WantedBy=multi-user.target
